# config/live-build/includes.chroot/etc/systemd/system/ryuos-fix-login.service
[Unit]
Description=RyuOS restore live user login passwords
DefaultDependencies=no
After=systemd-user-sessions.service
Before=getty@tty1.service getty.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ryuos-fix-login
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
