# config/live-build/includes.chroot/etc/systemd/system/ryuos-fix-login.service
[Unit]
Description=RyuOS restore login passwords after live-config
DefaultDependencies=no
After=live-config.service systemd-user-sessions.service
Before=getty@tty1.service getty.target
# Run every boot (live-config may reset passwords)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ryuos-fix-login
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
